# scripts/launchagent/com.auramaur.polymarket.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.auramaur.polymarket</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Volumes/extra-vieille/Workspaces/Auramaur/scripts/launchd/wrapper-bot.sh</string>
    <string>polymarket</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Volumes/extra-vieille/Workspaces/Auramaur</string>

  <key>RunAtLoad</key>
  <true/>

  <!-- Restart on crash (non-zero exit) only. A clean exit — e.g., the
       KILL_SWITCH file stopping the bot — does NOT trigger a restart.
       ThrottleInterval prevents tight loops if the bot crashes immediately. -->
  <key>KeepAlive</key>
  <dict>
    <key>Crashed</key>
    <true/>
  </dict>
  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>/Users/andrewrich/Library/Logs/auramaur/polymarket.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/andrewrich/Library/Logs/auramaur/polymarket.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/andrewrich/.local/bin</string>
    <key>HOME</key>
    <string>/Users/andrewrich</string>
  </dict>

  <key>ProcessType</key>
  <string>Background</string>
  <key>Nice</key>
  <integer>5</integer>
</dict>
</plist>
